# ipinfo.service
[Unit]
Description=ipinfo Flask app
After=network.target

[Service]
User=ipinfo
Group=ipinfo
Environment=BASE_DOMAIN=yourdomain.com
Environment=WHITELIST_IPS=
WorkingDirectory=/srv/ipinfo

# Setup runtime directory for the Unix socket
RuntimeDirectory=ipinfo
RuntimeDirectoryMode=0775

# Ensure logs directory exists
ExecStartPre=/bin/mkdir -p /srv/ipinfo/logs
ExecStartPre=/bin/chown ipinfo:ipinfo /srv/ipinfo/logs

# Set umask to allow group access to the socket
UMask=002

# Run Gunicorn with logging
ExecStart=/srv/ipinfo/venv/bin/gunicorn \
  --workers 3 \
  --bind unix:/run/ipinfo/ipinfo.sock \
  --access-logfile /srv/ipinfo/logs/access.log \
  --error-logfile /srv/ipinfo/logs/error.log \
  app:app

Restart=always

[Install]
WantedBy=multi-user.target
